# nginx-php7/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:php7.0-fpm]
command=/usr/sbin/php7.0-fpm
